# conf/solo-ci.service
[Unit]
Description = SmallSohoSolo Service

[Service]
ExecStart = /data/repo/solo_ci/solo-ci